# salt/roots/services/templates/upstart-service-template.conf
description     "{{ service_name }}"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec su - {{ run_as }} -c '{{ command_line }} >> {{ log_file }} 2>&1 '
